# session/supervisord.conf
# session/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:entrypoint]
command=/etc/entrypoint.sh
user=root
stdout_logfile=/tmp/entrypoint.log
redirect_stderr=true
autorestart=true
stopasgroup=true
priority=1

[program:pulseaudio]
command=bash -c "until [ -d \"%(ENV_XDG_RUNTIME_DIR)s\" ]; do sleep 0.5; done; /usr/bin/pulseaudio --daemonize=false --exit-idle-time=-1 --disallow-exit"
user=psp
environment=HOME="/home/psp",XDG_RUNTIME_DIR="%(ENV_XDG_RUNTIME_DIR)s",PULSE_RUNTIME_PATH="%(ENV_PULSE_RUNTIME_PATH)s"
stdout_logfile=/tmp/pulseaudio.log
redirect_stderr=true
autorestart=true
stopasgroup=true
priority=10

[program:selkies]
command=bash -c "until [ -S /tmp/.X11-unix/X20 ] && [ -S \"%(ENV_PULSE_RUNTIME_PATH)s/native\" ]; do sleep 0.5; done; /etc/selkies-gstreamer-entrypoint.sh"
user=psp
environment=HOME="/home/psp",USER="psp",SHELL="/bin/bash"
stdout_logfile=/tmp/selkies-gstreamer-entrypoint.log
redirect_stderr=true
autorestart=true
stopasgroup=true
stopsignal=INT
priority=20

[program:nginx]
command=bash -c "until nc -z localhost 8081; do sleep 0.5; done; /usr/sbin/nginx -g 'daemon off;'"
user=root
stdout_logfile=/tmp/nginx.log
redirect_stderr=true
autorestart=true
stopasgroup=true
priority=30

[program:game-supervisor]
command=bash -c "until [ -S /tmp/.X11-unix/X20 ]; do sleep 0.5; done; /usr/bin/python3 /opt/supervisor/supervisor.py"
user=psp
environment=HOME="/home/psp",DISPLAY="%(ENV_DISPLAY)s",XDG_RUNTIME_DIR="%(ENV_XDG_RUNTIME_DIR)s",PULSE_SERVER="%(ENV_PULSE_SERVER)s",SDL_JOYSTICK_DEVICE="/dev/input/js0"
stdout_logfile=/tmp/game-supervisor.log
redirect_stderr=true
autorestart=true
stopasgroup=true
priority=40
